# mactelnetd.service
[Unit]
Description=MAC-Telnet server
After=network.target

[Service]
ExecStart=/usr/sbin/mactelnetd -f

[Install]
WantedBy=multi-user.target
